# systemd-units/.config/systemd/user/kdeconnect-indicator.service
[Unit]
Description=KDEConnect Indicator
[Service]
ExecStart=/usr/bin/kdeconnect-indicator
[Install]
WantedBy=default.target
